# debian/postinst
#!/bin/bash -e

NAME=jmpapi
ETC=/etc/$NAME
RUN=/var/run/$NAME
LOG=/var/log/$NAME
CONFIG=$ETC/config.d
DOC=/usr/share/doc/$NAME

case "$1" in
    reconfigure|configure)
        # Install default config if one does not exist
        if [ ! -d $CONFIG ]; then
            mkdir -m 755 $CONFIG

            ln -s $DOC/config.d/jmpapi-auth.yaml $CONFIG/30-jmpapi-auth.yaml
            ln -s $DOC/config.d/jmpapi-docs.yaml $CONFIG/40-jmpapi-docs.yaml
            ln -s $DOC/config.d/local.yaml $CONFIG/50-local.yaml

            chmod a-r $CONFIG/50-local.yaml # May contain secrets
        fi

        # Setup permissions
        chown -R $NAME:root $ETC $RUN $LOG
        chmod -R u+r,a-w $ETC $RUN $LOG

        if [ -x insserv ]; then
            # Install the init.d script
            insserv -d $NAME

            # Restart the service
            service $NAME start || true

        else
            # Install the init.d script
            update-rc.d $NAME defaults 95

            # Start the service
            invoke-rc.d $NAME start || true
        fi
        ;;
esac
